# upstart/owncloud-sync.conf
description "file sync"
start on started unity8-dash 

exec /home/phablet/bin/owncloud-sync.sh
